# templates/default/god.conf.erb
#!upstart

description   "God Process Monitoring"
author        "Christoph Geschwind <christoph@mixxt.net>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn              # respawn the service if it dies
respawn limit 5 10   # stop respawning if it fails 5 times in 10 seconds

pre-start script
  mkdir -p /var/run/god  # create gods pid directory, you might not need this
end script

env PATH=/opt/rbenv/shims:/opt/rbenv/bin:/opt/rbenv/bin/rbenv:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
env RAILS_ENV=<%= node['cookbook_clarus']['rails_env'] %>
env RAILS_ROOT=<%= ::File.join(node["appbox"]["apps_dir"], node['cookbook_clarus']['appname'], 'current') %>
chdir <%= ::File.join(node["appbox"]["apps_dir"], node['cookbook_clarus']['appname'], 'current') %>
exec bundle exec god -D -c config/god/conf.god